# artifacts/service/com.nonaka.ai-harness.cliproxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MANAGEMENT_STATIC_PATH</key>
		<string>/Users/nonaka/.local/share/ai-harness/cliproxy/static</string>
		<key>WRITABLE_PATH</key>
		<string>/Users/nonaka/.local/share/ai-harness/cliproxy</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.nonaka.ai-harness.cliproxy</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/nonaka/.local/libexec/ai-harness/cli-proxy-api-start</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/nonaka/.local/share/ai-harness/cliproxy/logs/launchd.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/nonaka/.local/share/ai-harness/cliproxy/logs/launchd.out.log</string>
	<key>Umask</key>
	<integer>63</integer>
</dict>
</plist>
